# ansible/roles/minecraft/templates/minecraft.service.j2
[Unit]
Description=Minecraft Server
After=network.target

[Service]
Type=simple
User={{ minecraft_user }}
Group={{ minecraft_user }}
WorkingDirectory={{ minecraft_data_dir }}
ExecStart=/bin/bash -c 'SERVER_SH_INTERNAL=1 {{ minecraft_server_dir }}/server.sh start'
ExecStop=/bin/bash -c 'SERVER_SH_INTERNAL=1 {{ minecraft_server_dir }}/server.sh stop'
TimeoutStopSec=90
KillMode=mixed
KillSignal=SIGTERM
SuccessExitStatus=0 143
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=minecraft

# Limites système
LimitNOFILE=65536
MemoryMax={{ minecraft_memory_max }}

# Sécurité
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ minecraft_data_dir }} {{ minecraft_server_dir }}

[Install]
WantedBy=multi-user.target
